# nvda.service
[Unit]
Description=NVDA Stock Price Monitor
After=network-online.target

[Service]
WorkingDirectory=/opt/nvda_service
ExecStart=/opt/nvda_service/run.sh
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
